# docker/supervisord.conf
[unix_http_server]
file=/run/supervisord.sock

[supervisord]
logfile=/var/log/supervisord.log
loglevel=info
nodaemon=true

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///run/supervisord.sock

[include]
files = /etc/supervisor.d/*.ini

[program:laravel-worker-create-update-sql]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/artisan queue:work --queue=create-update-sql --tries=3
autostart=true
autorestart=true
numprocs=8
redirect_stderr=true
stdout_logfile=/var/www/storage/logs/worker.log

[program:laravel-worker-process-update-sql]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/artisan queue:work --queue=process-update-sql --tries=3
autostart=true
autorestart=true
numprocs=8
redirect_stderr=true
stdout_logfile=/var/www/storage/logs/worker.log
